# files/chezmoi/chezmoi-update.service
[Unit]
Description=Update Chezmoi dotfiles

[Service]
Type=oneshot
ExecStart=/usr/bin/chezmoi update --no-tty --keep-going
